# chef/local/james/files/default/mounta.upstart
description "mount all"

start on runlevel [2345]

exec mount -a

